# mmfw-sysconf-sc7730s/etc/profile.d/mmf.sh
export GST_DEBUG=2
export GST_REGISTRY=/home/app/.cache/gstreamer-1.0/registry.bin
export PULSE_CLIENTCONFIG=/etc/pulse/client.conf
